# hack-linux-for-v2f
#!/bin/bash

set -e

# you should be root in order to do these dangerous things
test "$(id -u)" = 0

# we want to have 0666 permission on all /dev/hidrawN devices automatically
echo 'KERNEL=="hidraw*", SUBSYSTEM=="hidraw", MODE="0666"' > /etc/udev/rules.d/99-hidrawN666.rules

# let the new rule takes effect now
udevadm trigger

# we want to have 0666 permission on /dev/uhid
( > /dev/uhid ; sleep 1 ; chmod 0666 /dev/uhid ) &

# we want to have 0666 permission on /dev/uhid even after rebooting
echo '#!/bin/bash
### BEGIN INIT INFO
# Provides: uhid666
# Required-Start:
# Required-Stop:
# Default-Start: 2
# Default-Stop:
### END INIT INFO
case "$1" in
    start)
        > /dev/uhid
        /bin/sleep 1
        /bin/chmod 0666 /dev/uhid
        exit 0
        ;;
    stop)
        exit 0
        ;;
    *)
        echo "Usage: $0 start|stop" >&2
        exit 3
        ;;
esac' > /etc/init.d/uhid666
chmod +x /etc/init.d/uhid666
update-rc.d uhid666 start 99 2 .
